# net-misc/pyload/files/pyload.initd
#!/sbin/runscript

NAME="PyLoad Daemon"

depend() {
	need net
}

start() {
	ebegin "Starting PyLoad"
		start-stop-daemon --start --quiet \
		--pidfile ${PIDFILE} \
		--exec /usr/bin/pyLoadCore -- --configdir=${CONFIGDIR} --daemon -p ${PIDFILE} ${OPTIONS}
	eend $?
}

stop() {
	ebegin "Stopping PyLoad"
		start-stop-daemon --stop --quiet \
		--pidfile ${PIDFILE}
	eend $?
}
